# nova.service
[Unit]
Description=Nova AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/nova
ExecStart=/usr/bin/python3 /home/pi/nova/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=/home/pi/nova/.env

[Install]
WantedBy=multi-user.target
